# build/supergiant-server.conf
# supergiant-server.conf
start on filesystem
script
  /etc/init.d/sgboot ; true
  eval $(cat /etc/environment | sed 's/^/export /')
  while ! echo exit | nc localhost 5432; do sleep 10; done
  /supergiant/generate_admin_user --config-file /etc/supergiant.json > /dev/ttyS0
  /supergiant/supergiant --config-file /etc/supergiant.json > /dev/ttyS0
end script
